# fully-automated-ops/systemd/net_guard.service
[Unit]
Description=Run net_guard probe/repair
[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
ExecStart=/usr/local/sbin/net_guard.sh
