# fully-automated-ops/systemd/panic_guard.service
[Unit]
Description=Run panic_guard (last-resort restart)
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/panic_guard.sh
