# fully-automated-ops/systemd/web_guard.service
[Unit]
Description=Run web_guard probe/repair
[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
ExecStart=/usr/local/sbin/web_guard.sh
